# deploy/agents.service
[Unit]
Description=Polymarket Betting Agent (NBA + NHL)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket-bot
ExecStartPre=/usr/bin/git -C /root/polymarket-bot fetch --all
ExecStartPre=/usr/bin/git -C /root/polymarket-bot reset --hard origin/master
ExecStart=/root/polymarket-bot/venv/bin/python orchestrator.py
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
