# deploy/agents.service
[Unit]
Description=Polymarket Betting Agent (NBA + NHL)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket-bot
ExecStart=/root/polymarket-bot/venv/bin/python orchestrator.py
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
