# deploy/agents.service
[Unit]
Description=NBA Polymarket Betting Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket-bot
ExecStart=/root/polymarket-bot/venv/bin/python -m nba_agent.main
Restart=always
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
