# deploy/polymarket-bot.service
[Unit]
Description=Polymarket 5-Min Trading Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarket-bot
ExecStart=/root/polymarket-bot/venv/bin/python /root/polymarket-bot/scalp_main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=append:/root/polymarket-bot/logs/bot.log
StandardError=append:/root/polymarket-bot/logs/bot.log

# Safety limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
